# config/cron/common/cron.daily/deploy
#!/bin/bash
MAILTO=""

/root/travaux/deploy.sh > /dev/null