# config/cron/web/cron.d/watchdog-mysql
# Monitoring MySQL
*/15 * * * * root /root/travaux/utils/watchdog_mysql.sh > /dev/null 2>&1
